# supervisord.conf


[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file


[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:gunicorn]
command=%(ENV_ENV_GUNICORN_CMD)s/gunicorn -w 3 memoflash_back.wsgi:application
directory=%(ENV_ENV_GUNICORN_DIR)s
user=root
autostart=true
autorestart=true
stderr_logfile=%(ENV_ENV_GUNICORN_DIR)s/gunicorn.err.log
stdout_logfile=%(ENV_ENV_GUNICORN_DIR)s/gunicorn.out.log
redirect_stderr=true
